# devguard.service
[Unit]
Description=DevGuard - LLM-Powered Security Scanning Agent
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/devguard
Environment="PATH=/home/ubuntu/devguard/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/ubuntu/devguard/venv/bin/streamlit run app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true --browser.gatherUsageStats=false
Restart=always
RestartSec=10

# Load environment variables from .env file
EnvironmentFile=/home/ubuntu/devguard/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=devguard

[Install]
WantedBy=multi-user.target

